# BigRedButton.conf
description "Big Red Button"

start on runlevel [2345]
chdir /opt/bigredbutton
exec /opt/bigredbutton/run.sh
